# settingsBackup/profile.d/java.sh
export JDK_DIR=/broncho/tools/jdk1.6.0_32
export JDK_HOME=$JDK_DIR
export PATH=$JDK_DIR/bin:$PATH
